# uwsgi.conf
# simple uWSGI script

description "orario uwsgi"
start on runlevel [2345]
stop on runlevel [06]

exec uwsgi /home/orario/orario/orario.ini